# systemd/waibar-desktop.service.in
# SPDX-FileCopyrightText: 2021-2025 Andrea Pappacoda <andrea@pappacoda.it>
# SPDX-License-Identifier: GPL-2.0-or-later

[Unit]
Description=@APPLICATION_NAME@ desktop user service
Wants=network-online.target

[Service]
ExecStart=@bindir@/@APPLICATION_EXECUTABLE@ --background
Restart=on-failure
RestartSec=5s
NoNewPrivileges=true
RestrictRealtime=true

[Install]
WantedBy=default.target
